# cdh/cdh.conf
[program:cloudera-scm-agent]
command=/bin/bash -c "service cloudera-scm-server start"

[program:dnsmasq]
command=/bin/bash -c "service dnsmasq start"
